# ansible/files/offload/offload_tidy.cron
# Managed by Ansible
0 4 * * * prom /usr/local/bin/offload_tidy | /usr/bin/logger -t offload_tidy
