# ansible/files/offload/tidy_rsync_logs.cron
# Managed by Ansible
0 3 * * * prom /bin/find /data/rsync_logs -name "rsync_*.log" -mtime +7 -delete > /dev/null
